# systemd/oci-log-rotation.service
[Unit]
Description=OCI log rotation and journald vacuum
After=systemd-journald.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/oci-log-rotation.sh
